# containers/dkfz_full/etc/profile.d/odcf-env.sh
#
# Copyright (c) 2023 German Cancer Research Center (DKFZ).
#
# Distributed under the MIT License (license terms are at https://github.com/TheRoddyWMS/AlignmentAndQCWorkflows).
#
export PATH=${PATH}:/software/bin
